# upstart_scripts/vrpn-rightglove.conf.in
# VRPN server
#

description	"VRPN device server for left serial glove"
start on ((local-filesystems and net-device-up IFACE!=lo)
          and stopped udevtrigger)
stop on runlevel [016]

respawn


env VRPN_HOME=@VRPN_HOME@
env VRPN_CFG_DIR=@VRPN_CFG_DIR@
env TAG=vrpn-rightglove

pre-start script
	logger -t "${TAG}" -- "Pre-start."
end script

script
	logger -t "${TAG}" -- "Device server starting up."
	${VRPN_HOME}/bin/up-progress
	export LD_LIBRARY_PATH=${VRPN_HOME}/lib

	# Based on http://serverfault.com/questions/114052/logging-a-daemons-output-with-upstart
	#mkfifo /tmp/${JOB}-fifo
	#(logger -t "${TAG}" < /tmp/${JOB}-fifo & )
	#exec > /tmp/${JOB}-fifo
	#rm /tmp/${JOB}-fifo

	exec ${VRPN_HOME}/bin/vrpn_server 3885 -f ${VRPN_CFG_DIR}/vrpn-rightserialglove.cfg
end script

pre-stop script
	logger -t "${TAG}" -- "Pre-stop."
end script

post-stop script
	${VRPN_HOME}/bin/down-progress
	logger -t "${TAG}" -- "Device server killed."
end script



